# files/usr/lib/systemd/user/obsidian-blue-hyprland-session.target
[Unit]
Description=obsidian-blue Hyprland graphical session
Documentation=https://github.com/afobsidian/obsidian-blue
Wants=graphical-session.target
BindsTo=graphical-session.target
Before=graphical-session.target
