# files/usr/lib/systemd/user/obsidian-blue-session.target
[Unit]
Description=obsidian-blue Hyprland session anchor
Wants=graphical-session.target
