# files/usr/lib/systemd/user/moniqued.service.d/10-obsidian-blue.conf
[Service]
ExecStartPre=/usr/libexec/obsidian-blue-configure-monique
